# linux/.config/autostart/albert.desktop
[Desktop Entry]
Type=Application
Exec=albert
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Albert
