# linux/.config/autostart/enpass.desktop
[Desktop Entry]
Type=Application
Exec=/opt/enpass/Enpass
Hidden=true
X-GNOME-Autostart-enabled=true
Name=Enpass

