# linux/.config/autostart/remap-keys.desktop
[Desktop Entry]
Type=Application
Exec=remap-keys
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Remap Keys

