# linux/.config/autostart/tilda.desktop
[Desktop Entry]
Type=Application
Exec=tilda
Hidden=false
X-GNOME-Autostart-enabled=true
Name=Tilda

