# osx/Library/LaunchAgents/my-jupyter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>my-jupyter</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/jupyter</string>
      <string>lab</string>
      <string>--port=8889</string>
      <string>--no-browser</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/spencerboucher</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>JUPYTER_DATA_DIR</key>
      <string>/Users/spencerboucher/.local/share/jupyter</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardErrorPath</key>
    <string>/Users/spencerboucher/Library/Logs/jupyter.err</string>
    <key>StandardOutPath</key>
    <string>/Users/spencerboucher/Library/Logs/jupyter.out</string>
  </dict>
</plist>
