# osx/Library/LaunchAgents/my-offlineimap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>KeepAlive</key>
    <true/>
    <key>Label</key>
    <string>my-offlineimap</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/zsh</string>
      <string>-c</string>
      <string>/usr/local/bin/offlineimap</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/offlineimap.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/offlineimap.out</string>
  </dict>
</plist>
