# linux/.config/systemd/user/jupyter.service
[Unit]
Description=Jupyter

[Service]
ExecStart=%h/.local/bin/jupyter-notebook --no-browser
Type=simple

[Install]
WantedBy=default.target
