# linux/.config/systemd/user/mailpile.service
[Unit]
Description=Mailpile

[Service]
ExecStart=/usr/bin/mailpile

[Install]
WantedBy=default.target
