# linux/.config/systemd/user/mailpile.service
[Unit]
Description=Mailpile

[Service]
ExecStart=/usr/bin/mailpile --www=localhost:33411 --wait
Type=simple

[Install]
WantedBy=default.target
