# gnome/.config/autostart/albert.desktop
[Desktop Entry]
Exec=albert
Hidden=false
Name=Albert
Type=Application
X-GNOME-Autostart-enabled=true
