# gnome/.config/autostart/albert.desktop
[Desktop Entry]
Exec=/usr/bin/albert
Hidden=false
Name=Albert
Type=Application
X-GNOME-Autostart-enabled=true
