# gnome/.config/autostart/remap-keys.desktop
[Desktop Entry]
Exec=remap-keys
Hidden=false
Name=Remap Keys
Type=Application
X-GNOME-Autostart-enabled=true
