# gnome/.config/autostart/remap-keys.desktop
[Desktop Entry]
Exec=/home/spencer/.local/bin/remap-keys
Hidden=false
Name=Remap Keys
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
