# gnome/.config/autostart/tilda.desktop
[Desktop Entry]
Exec=tilda
Hidden=false
Name=Tilda
Type=Application
X-GNOME-Autostart-enabled=true
